# scripts/linux/anotherclaw-backend.service.in
# AnotherClaw：systemd 单元模板（由 install-systemd-autostart.sh 生成并安装）
# Template for systemd — rendered by install-systemd-autostart.sh
#
# 行为：开机执行等价于「cd 仓库根 && just start-d」；关机/重启前「just close」。
# Boot: same as `cd <repo> && just start-d`; stop: `just close`.

[Unit]
Description=AnotherClaw backend (just start-d)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=@SVC_USER@
Group=@SVC_GROUP@
WorkingDirectory=@INSTALL_ROOT@
Environment=HOME=@SVC_HOME@
# 按需补充 just / podman / uv 所在路径（常见：/root/.cargo/bin、/root/.local/bin）
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:@EXTRA_PATH@
ExecStart=/bin/bash -lc 'just start-d'
ExecStop=/bin/bash -lc 'just close'
# 首次可能拉镜像、uv sync、构建前端，适当放宽
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
